# var/lib/dpkg/info/apt.conffiles
/etc/logrotate.d/apt
/etc/cron.daily/apt
/etc/apt/apt.conf.d/01autoremove
