# var/lib/dpkg/info/at.list
/.
/usr
/usr/share
/usr/share/doc/at
/usr/share/doc/at/changelog.gz
/usr/share/doc/at/timespec
/usr/share/doc/at/README
/usr/share/doc/at/Problems
/usr/share/doc/at/changelog.Debian.gz
/usr/share/doc/at/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/at.1.gz
/usr/share/man/man8
/usr/share/man/man8/atd.8.gz
/usr/share/man/man5
/usr/share/man/man5/at.allow.5.gz
/usr/bin
/usr/bin/at
/usr/bin/batch
/usr/sbin
/usr/sbin/atd
/etc
/etc/pam.d
/etc/pam.d/atd
/etc/at.deny
/etc/init.d
/etc/init.d/atd
/var
/var/spool
/var/spool/cron
/var/spool/cron/atspool
/var/spool/cron/atjobs
/usr/share/man/man1/batch.1.gz
/usr/share/man/man1/atrm.1.gz
/usr/share/man/man1/atq.1.gz
/usr/share/man/man5/at.deny.5.gz
/usr/bin/atrm
/usr/bin/atq
