# var/lib/dpkg/info/at.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    SEQFILE=/var/spool/cron/atjobs/.SEQ
    if [ ! -f $SEQFILE ]; then
	echo 0 > $SEQFILE
	chmod 600 $SEQFILE
	chown daemon:daemon $SEQFILE
    fi
fi

# Automatically added by dh_installinit
if [ -x "/etc/init.d/atd" ]; then
	update-rc.d atd defaults 89 11 >/dev/null
	invoke-rc.d atd start || exit $?
fi
# End automatically added section

