# var/lib/dpkg/info/bsdmainutils.list
/etc
/etc/cron.daily
/etc/cron.daily/bsdmainutils
/etc/default
/etc/default/bsdmainutils
/etc/calendar
/etc/calendar/default
