# var/lib/dpkg/info/exim4-base.conffiles
/etc/logrotate.d/exim4-base
/etc/logrotate.d/exim4-paniclog
/etc/cron.daily/exim4-base
/etc/init.d/exim4
