# var/lib/dpkg/info/exim4-base.list
/.
/etc
/etc/logrotate.d
/etc/logrotate.d/exim4-base
/etc/logrotate.d/exim4-paniclog
/etc/cron.daily
/etc/cron.daily/exim4-base
/etc/init.d
/etc/init.d/exim4
/usr
/usr/share
/usr/share/doc-base
/usr/share/doc-base/exim4-spec-txt
/usr/share/doc-base/exim4-filter-txt
/usr/share/doc-base/exim4-readme-debian
/usr/share/bug
/usr/share/bug/exim4-base
/usr/share/bug/exim4-base/script
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/exipick.8.gz
/usr/share/man/man8/exiqsumm.8.gz
/usr/share/man/man8/exinext.8.gz
/usr/share/man/man8/exiqgrep.8.gz
/usr/share/man/man8/exim_convert4r4.8.gz
/usr/share/man/man8/exigrep.8.gz
/usr/share/man/man8/exim_dbmbuild.8.gz
/usr/share/man/man8/exicyclog.8.gz
/usr/share/man/man8/eximstats.8.gz
/usr/share/man/man8/exim_lock.8.gz
/usr/share/man/man8/exim_db.8.gz
/usr/share/man/man8/exim_checkaccess.8.gz
/usr/share/man/man8/syslog2eximlog.8.gz
/usr/share/man/man8/exim.8.gz
/usr/share/man/man8/exiwhat.8.gz
/usr/share/doc/exim4-base
/usr/share/doc/exim4-base/README
/usr/share/doc/exim4-base/dbm.discuss.txt.gz
/usr/share/doc/exim4-base/filter.txt.gz
/usr/share/doc/exim4-base/changelog.Debian.gz
/usr/share/doc/exim4-base/README.UPDATING.gz
/usr/share/doc/exim4-base/Exim3.upgrade.gz
/usr/share/doc/exim4-base/NewStuff.gz
/usr/share/doc/exim4-base/README.Debian.gz
/usr/share/doc/exim4-base/changelog.gz
/usr/share/doc/exim4-base/copyright
/usr/share/doc/exim4-base/OptionLists.txt.gz
/usr/share/doc/exim4-base/README.Debian.html
/usr/share/doc/exim4-base/examples
/usr/share/doc/exim4-base/examples/logargs.sh
/usr/share/doc/exim4-base/examples/transport-filter.pl
/usr/share/doc/exim4-base/examples/aliases
/usr/share/doc/exim4-base/examples/unknownuser.sh
/usr/share/doc/exim4-base/examples/cramtest.pl
/usr/share/doc/exim4-base/examples/exim-gencert
/usr/share/doc/exim4-base/examples/exim-adduser
/usr/share/doc/exim4-base/examples/ratelimit.pl
/usr/share/doc/exim4-base/examples/example.conf.gz
/usr/share/doc/exim4-base/TODO.Debian
/usr/share/doc/exim4-base/Exim4.upgrade.gz
/usr/share/doc/exim4-base/NOTICE.gz
/usr/share/doc/exim4-base/spec.txt.gz
/usr/share/doc/exim4-base/NEWS.Debian.gz
/usr/share/doc/exim4-base/ACKNOWLEDGMENTS.gz
/usr/share/doc/exim4-base/GnuTLS-FAQ.txt.gz
/usr/share/doc/exim4-base/changelog.Debian.old.gz
/usr/share/doc/exim4-base/README.SIEVE.gz
/usr/sbin
/usr/sbin/exim_dbmbuild
/usr/sbin/exiwhat
/usr/sbin/exiqgrep
/usr/sbin/exipick
/usr/sbin/exim_checkaccess
/usr/sbin/exigrep
/usr/sbin/exim_dumpdb
/usr/sbin/exim_lock
/usr/sbin/syslog2eximlog
/usr/sbin/exim_tidydb
/usr/sbin/exiqsumm
/usr/sbin/eximstats
/usr/sbin/exim_convert4r4
/usr/sbin/exicyclog
/usr/sbin/exinext
/usr/sbin/exim_fixdb
/usr/share/man/man8/exim_dumpdb.8.gz
/usr/share/man/man8/exim_tidydb.8.gz
/usr/share/man/man8/exim_fixdb.8.gz
