# var/lib/dpkg/info/fake-hwclock.conffiles
/etc/default/fake-hwclock
/etc/init.d/fake-hwclock
/etc/cron.hourly/fake-hwclock
