# var/lib/dpkg/info/fake-hwclock.list
/.
/usr
/usr/share
/usr/share/doc/fake-hwclock
/usr/share/doc/fake-hwclock/changelog.gz
/usr/share/doc/fake-hwclock/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/fake-hwclock.8.gz
/sbin
/sbin/fake-hwclock
/etc
/etc/default
/etc/default/fake-hwclock
/etc/init.d
/etc/init.d/fake-hwclock
/etc/cron.hourly
/etc/cron.hourly/fake-hwclock
