# var/lib/dpkg/info/man-db.list
/var
/var/cache
/var/cache/man
/etc
/etc/cron.daily
/etc/cron.daily/man-db
/etc/cron.weekly
/etc/cron.weekly/man-db
/etc/manpath.config
