# var/lib/dpkg/info/minidlna.conffiles
/etc/init.d/minidlna
/etc/minidlna.conf
/etc/default/minidlna
