# var/lib/dpkg/info/minidlna.prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/minidlna" ]; then
	invoke-rc.d minidlna stop || exit $?
fi
# End automatically added section
