# var/lib/dpkg/info/minidlna.list
/.
/usr
/usr/bin
/usr/bin/minidlna
/usr/share
/usr/share/doc/minidlna
/usr/share/doc/minidlna/README
/usr/share/doc/minidlna/changelog.Debian.gz
/usr/share/doc/minidlna/TODO
/usr/share/doc/minidlna/changelog.gz
/usr/share/doc/minidlna/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/minidlna.1.gz
/usr/share/man/man5
/usr/share/man/man5/minidlna.conf.5.gz
/etc
/etc/init.d
/etc/init.d/minidlna
/etc/minidlna.conf
/etc/default
/etc/default/minidlna
