# var/lib/dpkg/info/mlocate.conffiles
/etc/updatedb.conf
/etc/cron.daily/mlocate
