# var/lib/dpkg/info/mpd.list
/.
/var
/var/run
/var/run/mpd
/var/log
/var/log/mpd
/var/lib
/var/lib/mpd
/var/lib/mpd/playlists
/var/lib/mpd/music
/usr
/usr/bin
/usr/bin/mpd
/usr/share
/usr/share/doc/mpd
/usr/share/doc/mpd/NEWS.gz
/usr/share/doc/mpd/README
/usr/share/doc/mpd/NEWS.Debian.gz
/usr/share/doc/mpd/README.Debian
/usr/share/doc/mpd/UPGRADING
/usr/share/doc/mpd/changelog.Debian.gz
/usr/share/doc/mpd/examples
/usr/share/doc/mpd/examples/mpd.conf.gz
/usr/share/doc/mpd/copyright
/usr/share/doc/mpd/AUTHORS
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/mpd
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mpd.1.gz
/usr/share/man/man5
/usr/share/man/man5/mpd.conf.5.gz
/etc
/etc/init.d
/etc/init.d/mpd
/etc/default
/etc/default/mpd
/etc/mpd.conf
/etc/logrotate.d
/etc/logrotate.d/mpd
