# var/lib/dpkg/info/mpd.conffiles
/etc/init.d/mpd
/etc/default/mpd
/etc/mpd.conf
/etc/logrotate.d/mpd
