# var/lib/dpkg/info/mutt.conffiles
/etc/Muttrc
/etc/Muttrc.d/charset.rc
/etc/Muttrc.d/colors.rc
/etc/Muttrc.d/compressed-folders.rc
/etc/Muttrc.d/gpg.rc
/etc/Muttrc.d/smime.rc
