# var/lib/dpkg/info/mutt.list
/.
/usr
/usr/lib
/usr/lib/mutt
/usr/lib/mutt/pgpewrap
/usr/lib/mutt/debian-ldap-query
/usr/lib/mutt/source-muttrc.d
/usr/lib/mutt/pgpring
/usr/lib/mutt/mailspell
/usr/lib/mutt/mailto-mutt
/usr/lib/mime
/usr/lib/mime/packages
/usr/lib/mime/packages/mutt
/usr/share
/usr/share/pixmaps
/usr/share/pixmaps/mutt.xpm
/usr/share/menu
/usr/share/menu/mutt
/usr/share/doc
/usr/share/doc/mutt
/usr/share/doc/mutt/TODO
/usr/share/doc/mutt/README.Debian
/usr/share/doc/mutt/changelog.Debian.gz
/usr/share/doc/mutt/NEWS.Debian.gz
/usr/share/doc/mutt/README
/usr/share/doc/mutt/PGP-Notes.txt.gz
/usr/share/doc/mutt/README.SSL.gz
/usr/share/doc/mutt/README.SMIME
/usr/share/doc/mutt/README.SECURITY
/usr/share/doc/mutt/README.Patches
/usr/share/doc/mutt/changelog.gz
/usr/share/doc/mutt/NEWS.old.gz
/usr/share/doc/mutt/copyright
/usr/share/doc/mutt/examples
/usr/share/doc/mutt/examples/Muttrc.gz
/usr/share/doc/mutt/examples/mutt_xtitle
/usr/share/doc/mutt/examples/Tin.rc
/usr/share/doc/mutt/examples/sample.muttrc-tlr.gz
/usr/share/doc/mutt/examples/pgp5.rc
/usr/share/doc/mutt/examples/colors.default
/usr/share/doc/mutt/examples/smime_keys_test.pl
/usr/share/doc/mutt/examples/sidebar.muttrc
/usr/share/doc/mutt/examples/sample.mailcap
/usr/share/doc/mutt/examples/pgp6.rc
/usr/share/doc/mutt/examples/pgp2.rc
/usr/share/doc/mutt/examples/Pine.rc
/usr/share/doc/mutt/examples/ray.muttrc
/usr/share/doc/mutt/examples/colors.angdraug.gz
/usr/share/doc/mutt/examples/sample.muttrc.gz
/usr/share/doc/mutt/examples/colors.linux
/usr/share/doc/mutt/examples/Mush.rc
/usr/share/doc/mutt/examples/gpg.rc
/usr/share/doc/mutt/examples/smime.rc
/usr/share/doc/mutt/NEWS.gz
/usr/share/doc/mutt/html
/usr/share/doc/mutt/html/optionalfeatures.html
/usr/share/doc/mutt/html/miscellany.html
/usr/share/doc/mutt/html/index.html
/usr/share/doc/mutt/html/manual.html
/usr/share/doc/mutt/html/security.html
/usr/share/doc/mutt/html/intro.html
/usr/share/doc/mutt/html/configuration.html
/usr/share/doc/mutt/html/advancedusage.html
/usr/share/doc/mutt/html/mimesupport.html
/usr/share/doc/mutt/html/tuning.html
/usr/share/doc/mutt/html/gettingstarted.html
/usr/share/doc/mutt/html/reference.html
/usr/share/doc/mutt/manual.txt.gz
/usr/share/bug
/usr/share/bug/mutt
/usr/share/bug/mutt/control
/usr/share/bug/mutt/script
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mutt.1.gz
/usr/share/man/man1/mutt_dotlock.1.gz
/usr/share/man/man1/smime_keys.1.gz
/usr/share/man/man5
/usr/share/man/man5/muttrc.5.gz
/usr/share/man/man5/mmdf.5.gz
/usr/share/man/man5/mbox.5.gz
/usr/share/applications
/usr/share/applications/mutt.desktop
/usr/share/locale
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/mutt.mo
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/uk/LC_MESSAGES/mutt.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/mutt.mo
/usr/share/locale/zh_TW
/usr/share/locale/zh_TW/LC_MESSAGES
/usr/share/locale/zh_TW/LC_MESSAGES/mutt.mo
/usr/share/locale/hu
/usr/share/locale/hu/LC_MESSAGES
/usr/share/locale/hu/LC_MESSAGES/mutt.mo
/usr/share/locale/el
/usr/share/locale/el/LC_MESSAGES
/usr/share/locale/el/LC_MESSAGES/mutt.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/mutt.mo
/usr/share/locale/et
/usr/share/locale/et/LC_MESSAGES
/usr/share/locale/et/LC_MESSAGES/mutt.mo
/usr/share/locale/eo
/usr/share/locale/eo/LC_MESSAGES
/usr/share/locale/eo/LC_MESSAGES/mutt.mo
/usr/share/locale/ca
/usr/share/locale/ca/LC_MESSAGES
/usr/share/locale/ca/LC_MESSAGES/mutt.mo
/usr/share/locale/eu
/usr/share/locale/eu/LC_MESSAGES
/usr/share/locale/eu/LC_MESSAGES/mutt.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/mutt.mo
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/mutt.mo
/usr/share/locale/tr
/usr/share/locale/tr/LC_MESSAGES
/usr/share/locale/tr/LC_MESSAGES/mutt.mo
/usr/share/locale/bg
/usr/share/locale/bg/LC_MESSAGES
/usr/share/locale/bg/LC_MESSAGES/mutt.mo
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/ja/LC_MESSAGES/mutt.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/mutt.mo
/usr/share/locale/ga
/usr/share/locale/ga/LC_MESSAGES
/usr/share/locale/ga/LC_MESSAGES/mutt.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/mutt.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/mutt.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/mutt.mo
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_MESSAGES
/usr/share/locale/zh_CN/LC_MESSAGES/mutt.mo
/usr/share/locale/gl
/usr/share/locale/gl/LC_MESSAGES
/usr/share/locale/gl/LC_MESSAGES/mutt.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/mutt.mo
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/id/LC_MESSAGES/mutt.mo
/usr/share/locale/sk
/usr/share/locale/sk/LC_MESSAGES
/usr/share/locale/sk/LC_MESSAGES/mutt.mo
/usr/share/locale/ko
/usr/share/locale/ko/LC_MESSAGES
/usr/share/locale/ko/LC_MESSAGES/mutt.mo
/usr/share/locale/lt
/usr/share/locale/lt/LC_MESSAGES
/usr/share/locale/lt/LC_MESSAGES/mutt.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/mutt.mo
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/mutt
/usr/share/doc-base
/usr/share/doc-base/mutt
/usr/bin
/usr/bin/smime_keys
/usr/bin/mutt_dotlock
/usr/bin/mutt-org
/etc
/etc/Muttrc
/etc/Muttrc.d
/etc/Muttrc.d/charset.rc
/etc/Muttrc.d/colors.rc
/etc/Muttrc.d/compressed-folders.rc
/etc/Muttrc.d/gpg.rc
/etc/Muttrc.d/smime.rc
