# var/lib/dpkg/info/nfs-common.list
/.
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/nfs-common
/usr/share/nfs-common
/usr/share/nfs-common/conffiles
/usr/share/nfs-common/conffiles/nfs-common.default
/usr/share/nfs-common/conffiles/idmapd.conf
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/nfsiostat.8.gz
/usr/share/man/man8/umount.nfs.8.gz
/usr/share/man/man8/nfsstat.8.gz
/usr/share/man/man8/showmount.8.gz
/usr/share/man/man8/gssd.8.gz
/usr/share/man/man8/rpcdebug.8.gz
/usr/share/man/man8/mount.nfs.8.gz
/usr/share/man/man8/blkmapd.8.gz
/usr/share/man/man8/sm-notify.8.gz
/usr/share/man/man8/idmapd.8.gz
/usr/share/man/man8/mountstats.8.gz
/usr/share/man/man8/statd.8.gz
/usr/share/man/man5
/usr/share/man/man5/nfs.5.gz
/usr/share/man/man5/nfsmount.conf.5.gz
/usr/share/man/man5/nfsidmap.5.gz
/usr/share/doc
/usr/share/doc/nfs-common
/usr/share/doc/nfs-common/changelog.gz
/usr/share/doc/nfs-common/changelog.Debian.gz
/usr/share/doc/nfs-common/README.Debian.nfsv4
/usr/share/doc/nfs-common/copyright
/usr/share/bug
/usr/share/bug/nfs-common
/usr/share/bug/nfs-common/control
/usr/share/bug/nfs-common/script
/usr/share/bug/nfs-utils
/usr/share/bug/nfs-utils/control
/usr/share/bug/nfs-utils/presubj
/usr/sbin
/usr/sbin/mountstats
/usr/sbin/blkmapd
/usr/sbin/rpc.svcgssd
/usr/sbin/rpc.idmapd
/usr/sbin/nfsidmap
/usr/sbin/rpcdebug
/usr/sbin/rpc.gssd
/usr/sbin/start-statd
/usr/sbin/nfsiostat
/usr/sbin/gss_clnt_send_err
/usr/sbin/nfsstat
/usr/sbin/gss_destroy_creds
/var
/var/lib
/var/lib/nfs
/var/lib/nfs/sm.bak
/var/lib/nfs/sm
/var/lib/nfs/rpc_pipefs
/var/lib/nfs/state
/sbin
/sbin/rpc.statd
/sbin/sm-notify
/sbin/mount.nfs
/sbin/showmount
/sbin/osd_login
/etc
/etc/init.d
/etc/init.d/nfs-common
/usr/share/man/man8/rpc.sm-notify.8.gz
/usr/share/man/man8/rpc.idmapd.8.gz
/usr/share/man/man8/rpc.gssd.8.gz
/usr/share/man/man8/rpc.statd.8.gz
/sbin/umount.nfs4
/sbin/mount.nfs4
/sbin/umount.nfs
