# var/lib/dpkg/info/nfs-common.prerm
#!/bin/sh

set -e



case "$1" in
    remove|purge)
	[ -x /etc/init.d/nfs-common ] &&
	    invoke-rc.d nfs-common stop
	;;
esac

if [ "$1" != upgrade ]
then
    rm -f /var/lib/nfs/sm/*      \
          /var/lib/nfs/sm.bak/*  \
          /var/lib/nfs/state
fi

exit 0
