# var/lib/dpkg/info/nginx-common.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/nginx-common
/usr/share/doc/nginx-common/copyright
/usr/share/doc/nginx-common/changelog.gz
/usr/share/doc/nginx-common/changelog.Debian.gz
/usr/share/nginx
/usr/share/nginx/www
/usr/share/nginx/www/index.html
/usr/share/nginx/www/50x.html
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/nginx.1.gz
/etc
/etc/nginx
/etc/nginx/win-utf
/etc/nginx/koi-win
/etc/nginx/naxsi.rules
/etc/nginx/uwsgi_params
/etc/nginx/naxsi_core.rules
/etc/nginx/sites-enabled
/etc/nginx/scgi_params
/etc/nginx/mime.types
/etc/nginx/proxy_params
/etc/nginx/sites-available
/etc/nginx/sites-available/default
/etc/nginx/koi-utf
/etc/nginx/nginx.conf
/etc/nginx/fastcgi_params
/etc/nginx/conf.d
/etc/init.d
/etc/init.d/nginx
/etc/default
/etc/default/nginx
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/nginx
/etc/logrotate.d
/etc/logrotate.d/nginx
/var
/var/log
/var/log/nginx
/var/lib
/var/lib/nginx
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/nginx.service
/usr/share/doc/nginx-common/CHANGES.gz
