# var/lib/dpkg/info/ntp.list
/.
/usr
/usr/sbin
/usr/sbin/ntpd
/usr/sbin/ntp-keygen
/usr/sbin/ntp-wait
/usr/sbin/ntptime
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ntp-wait.8.gz
/usr/share/man/man8/ntpd.8.gz
/usr/share/man/man8/ntptime.8.gz
/usr/share/man/man8/ntp-keygen.8.gz
/usr/share/man/man1
/usr/share/man/man1/ntptrace.1.gz
/usr/share/man/man1/ntpq.1.gz
/usr/share/man/man1/ntpdc.1.gz
/usr/share/man/man1/ntpsweep.1.gz
/usr/share/man/man1/sntp.1.gz
/usr/share/man/man5
/usr/share/man/man5/ntp.conf.5.gz
/usr/share/doc/ntp
/usr/share/doc/ntp/changelog.Debian.gz
/usr/share/doc/ntp/NEWS.Debian.gz
/usr/share/doc/ntp/copyright
/usr/share/doc/ntp/NEWS.gz
/usr/share/doc/ntp/README.versions
/usr/share/doc/ntp/README.refclocks
/usr/share/doc/ntp/README.Debian.gz
/usr/bin
/usr/bin/ntpsweep
/usr/bin/sntp
/usr/bin/ntpq
/usr/bin/ntpdc
/usr/bin/ntptrace
/etc
/etc/dhcp
/etc/dhcp/dhclient-exit-hooks.d
/etc/dhcp/dhclient-exit-hooks.d/ntp
/etc/ntp.conf
/etc/init.d
/etc/init.d/ntp
/etc/default
/etc/default/ntp
/etc/cron.daily
/etc/cron.daily/ntp
/var
/var/lib
/var/lib/ntp
/var/log
/var/log/ntpstats
