# var/lib/dpkg/info/openssh-server.list
/.
/usr
/usr/sbin
/usr/sbin/sshd
/usr/lib
/usr/lib/openssh
/usr/lib/openssh/sftp-server
/usr/share
/usr/share/doc/openssh-client
/usr/share/doc/openssh-client/examples
/usr/share/doc/openssh-client/examples/sshd_config
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/sshd_config.5.gz
/usr/share/man/man8
/usr/share/man/man8/sshd.8.gz
/usr/share/man/man8/sftp-server.8.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/openssh-server
/etc
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/openssh-server
/etc/pam.d
/etc/pam.d/sshd
/etc/default
/etc/default/ssh
/etc/init.d
/etc/init.d/ssh
/etc/network
/etc/network/if-up.d
/etc/network/if-up.d/openssh-server
/usr/lib/sftp-server
/usr/share/doc/openssh-server
/usr/share/man/man5/authorized_keys.5.gz
