# var/lib/dpkg/info/php5-fpm.list
/.
/usr
/usr/sbin
/usr/sbin/php5-fpm
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/php5-fpm
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/php5-fpm.8.gz
/usr/share/doc
/usr/lib
/usr/lib/php5
/usr/lib/php5/20100525+lfs
/etc
/etc/logrotate.d
/etc/logrotate.d/php5-fpm
/etc/init.d
/etc/init.d/php5-fpm
/etc/php5
/etc/php5/fpm
/etc/php5/fpm/pool.d
/etc/php5/fpm/pool.d/www.conf
/etc/php5/fpm/php-fpm.conf
/usr/share/doc/php5-fpm
/etc/php5/fpm/conf.d
