# var/lib/dpkg/info/php5-common.list
/.
/usr
/usr/sbin
/usr/sbin/php5enmod
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/php5-common
/usr/share/doc
/usr/share/doc/php5-common
/usr/share/doc/php5-common/changelog.Debian.gz
/usr/share/doc/php5-common/NEWS.Debian.gz
/usr/share/doc/php5-common/copyright
/usr/share/doc/php5-common/EXTENSIONS.gz
/usr/share/doc/php5-common/README.Debian.security
/usr/share/doc/php5-common/changelog.gz
/usr/share/doc/php5-common/CREDITS
/usr/share/doc/php5-common/CODING_STANDARDS.gz
/usr/share/doc/php5-common/README.PHP4-TO-PHP5-THIN-CHANGES.gz
/usr/share/doc/php5-common/README.EXT_SKEL.gz
/usr/share/doc/php5-common/test-results.txt.gz
/usr/share/doc/php5-common/README.Debian.gz
/usr/share/doc/php5-common/README.SELF-CONTAINED-EXTENSIONS.gz
/usr/share/php5
/usr/share/php5/common
/usr/share/php5/common/pdo.ini
/usr/share/php5/php.ini-development
/usr/share/php5/php.ini-production.cli
/usr/share/php5/php.ini-production
/usr/lib
/usr/lib/php5
/usr/lib/php5/libexec
/usr/lib/php5/20100525+lfs
/usr/lib/php5/20100525+lfs/pdo.so
/usr/lib/php5/maxlifetime
/etc
/etc/cron.d
/etc/cron.d/php5
/etc/php5
/etc/php5/conf.d
/etc/php5/mods-available
/var
/var/lib
/var/lib/php5
/usr/sbin/php5dismod
