# var/lib/dpkg/info/plymouth.list
/.
/sbin
/sbin/plymouthd
/lib
/lib/arm-linux-gnueabihf
/lib/arm-linux-gnueabihf/libply.so.2.1.0
/lib/arm-linux-gnueabihf/libply-splash-core.so.2.1.0
/lib/systemd
/lib/systemd/system
/lib/systemd/system/plymouth-read-write.service
/lib/systemd/system/plymouth-reboot.service
/lib/systemd/system/plymouth-start.service
/lib/systemd/system/plymouth-quit.service
/lib/systemd/system/plymouth-poweroff.service
/lib/systemd/system/plymouth-kexec.service
/lib/systemd/system/plymouth-halt.service
/lib/systemd/system/plymouth-quit-wait.service
/bin
/bin/plymouth
/etc
/etc/plymouth
/etc/plymouth/plymouthd.conf
/etc/init.d
/etc/init.d/plymouth-log
/etc/init.d/plymouth
/usr
/usr/sbin
/usr/sbin/plymouth-set-default-theme
/usr/share
/usr/share/doc/plymouth
/usr/share/doc/plymouth/changelog.Debian.gz
/usr/share/doc/plymouth/README.Debian
/usr/share/doc/plymouth/copyright
/usr/share/doc/plymouth/changelog.gz
/usr/share/doc/plymouth/README
/usr/share/doc/plymouth/TODO
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/plymouth
/usr/share/initramfs-tools/scripts/init-premount
/usr/share/initramfs-tools/scripts/init-premount/plymouth
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/plymouth
/usr/share/plymouth
/usr/share/plymouth/plymouthd.defaults
/usr/share/plymouth/debian-logo.png
/usr/share/plymouth/themes
/usr/share/plymouth/themes/text
/usr/share/plymouth/themes/text/text.plymouth
/usr/share/plymouth/themes/details
/usr/share/plymouth/themes/details/details.plymouth
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/plymouth
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/plymouth.8.gz
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/plymouth
/usr/lib/arm-linux-gnueabihf/plymouth/text.so
/usr/lib/arm-linux-gnueabihf/plymouth/plymouth-generate-initrd
/usr/lib/arm-linux-gnueabihf/plymouth/plymouth-update-initrd
/usr/lib/arm-linux-gnueabihf/plymouth/plymouth-populate-initrd
/usr/lib/arm-linux-gnueabihf/plymouth/details.so
/usr/lib/arm-linux-gnueabihf/libply-boot-client.so.2.1.0
/lib/arm-linux-gnueabihf/libply.so.2
/lib/arm-linux-gnueabihf/libply-splash-core.so.2
/usr/share/man/man8/plymouth-set-default-theme.8.gz
/usr/lib/arm-linux-gnueabihf/libply-boot-client.so.2
