# var/lib/dpkg/info/rpcbind.list
/.
/etc
/etc/insserv.conf.d
/etc/insserv.conf.d/rpcbind
/etc/init.d
/etc/init.d/rpcbind
/usr
/usr/share
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/rpcinfo.7.gz
/usr/share/man/man8
/usr/share/man/man8/rpcbind.8.gz
/usr/share/doc/rpcbind
/usr/share/doc/rpcbind/changelog.Debian.gz
/usr/share/doc/rpcbind/changelog.gz
/usr/share/doc/rpcbind/copyright
/usr/sbin
/usr/sbin/rpcinfo
/sbin
/sbin/rpcbind
