# var/lib/dpkg/info/rsync.list
/etc
/etc/default
/etc/default/rsync
/etc/init.d
/etc/init.d/rsync
