# var/lib/dpkg/info/rsync.postrm
#!/bin/sh -e
# postrm for rsync

if [ "$1" = "purge" ] ; then
    update-rc.d rsync remove >/dev/null
fi

exit 0
