# var/lib/dpkg/info/samba.list
/.
/usr
/usr/lib
/usr/lib/samba
/usr/lib/samba/vfs
/usr/lib/samba/vfs/shadow_copy2.so
/usr/lib/samba/vfs/fileid.so
/usr/lib/samba/vfs/audit.so
/usr/lib/samba/vfs/netatalk.so
/usr/lib/samba/vfs/default_quota.so
/usr/lib/samba/vfs/preopen.so
/usr/lib/samba/vfs/xattr_tdb.so
/usr/lib/samba/vfs/smb_traffic_analyzer.so
/usr/lib/samba/vfs/expand_msdfs.so
/usr/lib/samba/vfs/cap.so
/usr/lib/samba/vfs/streams_xattr.so
/usr/lib/samba/vfs/time_audit.so
/usr/lib/samba/vfs/acl_xattr.so
/usr/lib/samba/vfs/acl_tdb.so
/usr/lib/samba/vfs/syncops.so
/usr/lib/samba/vfs/linux_xfs_sgid.so
/usr/lib/samba/vfs/extd_audit.so
/usr/lib/samba/vfs/crossrename.so
/usr/lib/samba/vfs/catia.so
/usr/lib/samba/vfs/full_audit.so
/usr/lib/samba/vfs/readahead.so
/usr/lib/samba/vfs/readonly.so
/usr/lib/samba/vfs/fake_perms.so
/usr/lib/samba/vfs/recycle.so
/usr/lib/samba/vfs/dirsort.so
/usr/lib/samba/vfs/scannedonly.so
/usr/lib/samba/vfs/shadow_copy.so
/usr/lib/samba/vfs/streams_depot.so
/usr/share
/usr/share/doc
/usr/share/doc/samba
/usr/share/doc/samba/examples
/usr/share/doc/samba/examples/LDAP
/usr/share/doc/samba/examples/LDAP/get_next_oid
/usr/share/doc/samba/examples/LDAP/samba.schema.oc.IBM-DS
/usr/share/doc/samba/examples/LDAP/samba-schema-FDS.ldif.gz
/usr/share/doc/samba/examples/LDAP/samba-nds.schema.gz
/usr/share/doc/samba/examples/LDAP/samba.schema.gz
/usr/share/doc/samba/examples/LDAP/samba-schema-netscapeds5.x.README
/usr/share/doc/samba/examples/LDAP/README
/usr/share/doc/samba/examples/LDAP/samba.schema.at.IBM-DS.gz
/usr/share/doc/samba/examples/LDAP/samba.ldif.gz
/usr/share/doc/samba/examples/LDAP/ol-schema-migrate.pl.gz
/usr/share/doc/samba/examples/LDAP/samba-schema.IBMSecureWay
/usr/share/doc/samba/examples/LDAP/convertSambaAccount.gz
/usr/share/doc/samba/changelog.Debian.gz
/usr/share/doc/samba/NEWS.Debian.gz
/usr/share/doc/samba/diagnosis.html
/usr/share/doc/samba/copyright
/usr/share/doc/samba/README.Debian
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/smbcontrol.1.gz
/usr/share/man/man1/smbstatus.samba3.1.gz
/usr/share/man/man1/profiles.1.gz
/usr/share/man/man8
/usr/share/man/man8/smbd.8.gz
/usr/share/man/man8/mksmbpasswd.8.gz
/usr/share/man/man8/nmbd.8.gz
/usr/share/man/man8/eventlogadm.8.gz
/usr/share/man/man8/pdbedit.8.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/samba
/usr/bin
/usr/bin/smbstatus.samba3
/usr/bin/smbcontrol
/usr/bin/pdbedit
/usr/bin/profiles
/usr/bin/eventlogadm
/usr/sbin
/usr/sbin/nmbd
/usr/sbin/mksmbpasswd
/usr/sbin/smbd
/etc
/etc/cron.daily
/etc/cron.daily/samba
/etc/init.d
/etc/init.d/samba
/etc/logrotate.d
/etc/logrotate.d/samba
/var
/var/lib
/var/lib/samba
/var/lib/samba/printers
/var/lib/samba/printers/W32MIPS
/var/lib/samba/printers/WIN40
/var/lib/samba/printers/W32ALPHA
/var/lib/samba/printers/W32PPC
/var/lib/samba/printers/W32X86
/var/lib/samba/printers/COLOR
/var/lib/samba/printers/x64
/var/lib/samba/printers/IA64
/var/spool
/var/spool/samba
