# var/lib/dpkg/info/triggerhappy.list
/.
/usr
/usr/sbin
/usr/sbin/thd
/usr/sbin/th-cmd
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/thd.1.gz
/usr/share/man/man1/th-cmd.1.gz
/usr/share/doc
/usr/share/doc/triggerhappy
/usr/share/doc/triggerhappy/changelog.Debian.gz
/usr/share/doc/triggerhappy/copyright
/usr/share/doc/triggerhappy/README
/usr/share/doc/triggerhappy/triggerhappy.conf.examples
/etc
/etc/triggerhappy
/etc/triggerhappy/triggers.d
/etc/init.d
/etc/init.d/triggerhappy
/etc/default
/etc/default/triggerhappy
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-triggerhappy.rules
