# var/lib/dpkg/info/triggerhappy.conffiles
/etc/init.d/triggerhappy
/etc/default/triggerhappy
