# var/lib/dpkg/info/wicd-daemon.list
/etc
/etc/init.d
/etc/init.d/wicd
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/wicd.conf
/etc/default
/etc/default/wicd
/etc/wicd
/etc/wicd/dhclient.conf.template.default
/etc/wicd/encryption
/etc/wicd/encryption/templates
/etc/wicd/encryption/templates/active
/etc/wicd/encryption/templates/peap-tkip
/etc/wicd/encryption/templates/wpa-psk
/etc/wicd/encryption/templates/wpa2-peap
/etc/wicd/encryption/templates/active_wired
/etc/wicd/encryption/templates/ttls
/etc/wicd/encryption/templates/wpa-peap
/etc/wicd/encryption/templates/eap-tls
/etc/wicd/encryption/templates/peap
/etc/wicd/encryption/templates/wpa2-leap
/etc/wicd/encryption/templates/wep-hex
/etc/wicd/encryption/templates/eap
/etc/wicd/encryption/templates/wpa
/etc/wicd/encryption/templates/wep-passphrase
/etc/wicd/encryption/templates/leap
/etc/wicd/encryption/templates/wired_8021x
/etc/wicd/encryption/templates/psu
/etc/wicd/encryption/templates/wep-shared
/var
/var/log
/var/log/wicd
/var/lib
/var/lib/wicd
